# deploy/lightsail/tr-signal.service
[Unit]
Description=terminal-relay signaling server
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/tr-signal --addr :80 --webroot /opt/tr-web
Restart=always
RestartSec=2
User=trsignal
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
